# build/supervisord.conf
# https://advancedweb.hu/supervisor-with-docker-lessons-learned/

# To control the state of programs, use supervisorctl:
# supervisorctl start <program>
# supervisorctl stop <program>
# supervisorctl restart <program>
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/sbin/php-fpm7.3 -F
autorestart=true
# Logging
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
numprocs=1
startsecs=0
process_name=%(program_name)s_%(process_num)02d
killasgroup=true
stopasgroup=true
# Logging
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
